# hydra_scripts/initrd/etc/inittab
# /etc/inittab

::sysinit:/linuxrc

::respawn:-/bin/sh
#ttyS1::respawn:-/bin/sh
